# backend/deploy/com.dragonreserve.backend.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.dragonreserve.backend</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Volumes/disk/project/dragonReserveSystem/backend/.venv/bin/python</string>
    <string>-m</string>
    <string>uvicorn</string>
    <string>app.main:app</string>
    <string>--host</string>
    <string>127.0.0.1</string>
    <string>--port</string>
    <string>8000</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Volumes/disk/project/dragonReserveSystem/backend</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Volumes/disk/project/dragonReserveSystem/backend/logs/backend.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Volumes/disk/project/dragonReserveSystem/backend/logs/backend.stderr.log</string>
</dict>
</plist>
